# prepare_oim/roles/deploy_containers/common/templates/omnia.service.j2
[Unit]
Description=Top-level target for Omnia Core and OpenCHAMI
Requires=omnia_core.service openchami.target pulp.service registry.service minio.service {{ auth_service }} {{ kubespray_service_list }}
After=network.target
Wants=network-online.target

[Install]
WantedBy=multi-user.target default.target
