# prepare_oim/roles/deploy_containers/common/templates/omnia.service.j2
[Unit]
Description=Top-level target for Omnia Core and OpenCHAMI
Requires=omnia_core.service openchami.target
After=network.target
Wants=network-online.target
 
[Install]
WantedBy=multi-user.target